# motel/systemd/portal-qa.timer
[Unit]
Description=Run West Bethel Portal QA every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Unit=portal-qa.service
Persistent=true

[Install]
WantedBy=timers.target
